# install/virtex.service
[Unit]
Description=Start Virtex USB-HID System
After=network.target

[Service]
ExecStart=/usr/bin/virtex-start-usb
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
